# images/ubuntu/files/gha-runner.service
[Unit]
Description=GitHub Actions Runner
After=network.target

[Service]
ExecStart=/usr/local/bin/start-gha-runner
User=gha
Group=gha
WorkingDirectory=/gha

[Install]
WantedBy=multi-user.target
